# infra/rankup/boringd-isolated-worker.service
[Unit]
Description=boringd isolated Firecracker build worker
After=network-online.target boring-net.service
Wants=network-online.target
Requires=boring-net.service

[Service]
Type=simple
EnvironmentFile=/etc/boring/isolated-worker.env
LoadCredential=boring_token:/etc/boring/boring_token
Environment=BORING_TOKEN_FILE=%d/boring_token
ExecStartPre=/usr/local/bin/boringd -check-config
ExecStartPre=/opt/boring/bin/net-setup.sh --check
ExecStart=/usr/local/bin/boringd
WorkingDirectory=/opt/boring
Restart=on-failure
RestartSec=2
User=root
Group=root
UMask=0077
Delegate=yes
LimitNOFILE=65536
LimitNPROC=4096
RuntimeDirectory=boring
StateDirectory=boring
PrivateTmp=yes
ProtectHome=yes
ProtectHostname=yes
ProtectClock=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
